# videoserver/rootfs/etc/services.d/videoserver/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start videoserver service
# ==============================================================================
function startStream()
{
	while [ 1 ] ; do
		sleep 1s
		echo "start ffmpeg $1 from input: $2"
		#pid=$(ps aux | grep -v grep | grep ffmpeg | grep "$1" | grep "$2")
		#if [ -n "$pid"  ] ; then
		#	kill -9 $(echo $pid | cut -d ' ' -f1)
		#fi
		ffmpeg -i "$2" -strict -2 "http://localhost:8090/$1.ffm"
		echo "ffmpeg exit done errorcode: $?"
	done
}

function parseFFMpegConf()
{
	json=/data/options.json
	for row in $(jq -cr '.video_sources[]' $json); do
		name=$(echo $row | jq '.name')
		input=$(echo $row | jq '.input')
		if [ -n "$name" && -n "$input" ] ; then
			startStream "$name" "$input" &
		else
			echo "Error invalid config, misssing name or input: $row"
		fi
	done
}

parseFFMpegConf

echo "starting ffserver ..."
while [ 1 ] ; do
	ffserver
	sleep 1s
done
