# videoserver/rootfs/etc/services.d/videoserver/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start videoserver service
# ==============================================================================
echo "starting ffserver ..."
ls /tmp/ffmpeg
while [ 1 ] ; do
	ffserver
	sleep 1s
done
