# videoserver/rootfs/etc/services.d/videoserver/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start videoserver service
# ==============================================================================

echo "Starting ffserver"

ffserver

echo "ffserver exited, errorcode: $?"
